# example/sdnotify-example.service.template
[Unit]
Description=Example Service for Java SDNotify
After=network.target

[Service]
Type=notify-reload
ExecStart=${appPath}
WatchdogSec=5
Restart=always

[Install]
WantedBy=multi-user.target